# modules/whimsy_server/files/whimsy-pubsub.conf
description "Watches for and deploys updates to the whimsy repository"

start on runlevel [2345]

respawn

exec /usr/local/bin/ruby2.3.1 /srv/whimsy/tools/pubsub.rb --puppet
